# mitmproxy/supervisord.conf
[program:mitmproxy]
process_name = "openvpn"
autorestart = true
redirect_stderr=true
command=mitmweb --no-web-open-browser  --web-host 0.0.0.0  --web-port 8080]
