# burp/supervisord.conf
[supervisord]
nodaemon=true

[program:burp]
autorestart = true
command =

[program:mitmweb]
autorestart = true
command = mitmweb --no-web-open-browser  --web-host 0.0.0.0 --web-port 8080

